# root/etc/services.d/twtxt/run
#!/usr/bin/with-contenv bash

cd /app || exit

exec \
	s6-setuidgid abc /app/twtd \
	  -d /config \
    -s bitcask:///data/twtxt.db \
    -u http://127.0.0.1:8000
